# systemd/findsenryu4discord.target
[Unit]
Description=FindSenryu4Discord stack

[Install]
WantedBy=multi-user.target
